# roles/cmd_center/templates/ansible-proxmox.timer.j2
[Unit]
Description=Run Ansible infrastructure playbooks every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
