# roles/cmd_center/templates/ansible-proxmox.timer.j2
[Unit]
Description=Run Ansible Proxmox playbook every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
