# roles/monitoring/pve_quorum/templates/pve-quorum-metrics.service.j2
[Unit]
Description=Collect PVE quorum metrics for node_exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pve-quorum-metrics.sh
