# roles/pve/hookscripts/templates/gpu-vm-monitor.service.j2
[Unit]
Description=GPU VM Monitor — detects and recovers from GPU failures via PCI reset

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpu-vm-monitor.sh
