# roles/pve/hookscripts/templates/gpu-vm-monitor.timer.j2
[Unit]
Description=GPU VM Monitor timer — checks GPU health every 30s

[Timer]
OnBootSec=60
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
