# roles/jellyfin/templates/gpu-watchdog.service.j2
[Unit]
Description=GPU Watchdog — auto-recovers from GPU failures
After=jellyfin.service

[Service]
Type=simple
ExecStart=/usr/local/bin/gpu-watchdog.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
